# files/etc/rc.local
echo "AT+CFUN=1" > /dev/ttyUSB3
exit 0

